# example.service
[Unit]
Description=IP Synchronisation service

[Service]
Type=idle
ExecStart=/path/to/run_sync_service.sh

[Install]
WantedBy=multi-user.target
